# server/aesdsocket-start-stop
#!/bin/sh
### BEGIN INIT INFO
# Provides:          aesdsocket
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Start/stop aesdsocket daemon
### END INIT INFO

DAEMON_PATH="/usr/bin/aesdsocket"
PIDFILE="/var/run/aesdsocket.pid"

case "$1" in
    start)
        echo "Starting aesdsocket..."
        start-stop-daemon --start --background --make-pidfile --pidfile "$PIDFILE" --exec "$DAEMON_PATH" -- -d
        ;;
    stop)
        echo "Stopping aesdsocket..."
        start-stop-daemon --stop --pidfile "$PIDFILE" --retry TERM/5
        rm -f "$PIDFILE"
        ;;
    restart)
        $0 stop
        $0 start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
        ;;
esac

exit 0
